# systemd/stringcast-wayland-listener.service
[Unit]
Description=Stringcast Wayland Listener (evdev/uinput)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# PoC default: root can read /dev/input/event* and write to uinput.
# For a lower-privilege setup, create a dedicated user and apply the sample
# udev rules in udev/99-stringcast-input.rules.
User=root
Group=root
WorkingDirectory=/opt/stringcast
EnvironmentFile=-/etc/stringcast/wayland-listener.env
ExecStart=/opt/stringcast/.venv/bin/python /opt/stringcast/scripts/wayland_listener.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
